# infra/node-exporter/node-exporter.service
[Unit]
Description=node-exporter

[Service]
ExecStart=/home/ubuntu/node_exporter/node_exporter

[Install]
WantedBy=multi-user.target